# packages/hardware-monitoring/watchdog/files/systemd/watchdog.service
[Unit]
Description=watchdog daemon

[Service]
ExecStart=/usr/sbin/watchdog -F

[Install]
WantedBy=multi-user.target
