# packages/hardware-monitoring/mcelog/files/systemd/mcelog.service
[Unit]
Description=mcelog daemon

[Service]
ExecStart=/usr/sbin/mcelog --daemon --foreground

[Install]
WantedBy=multi-user.target
